# roles/monitoring/templates/node_exporter.service.j2
[Unit]
Description=Node Exporter

[Service]
User=root
ExecStart=/bin/node_exporter

[Install]
WantedBy=default.target
